# cubic/project/rootfs-overlay/etc/systemd/system/reswarmos-update-install.service

[Unit]
Description=Install ReswarmOS update bundle

[Service]
Type=oneshot
ExecStart=rauc install $(ls -t /tmp/*.raucb)
ExecStop=rm /tmp/*.raucb -v

[Install]

