# reswarmify/rootfs/etc/systemd/system/dbus.service.d/restart.conf
[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
