# buildroot/rootfs/etc/systemd/system/reagent-upgrade.path
[Unit]
Description=Reagent upgrade watcher
ConditionPathExists=/opt/reagent/reswarm-mode
#After=reswarm.service network-online.target
#Requires=network-online.target

[Path]
PathExistsGlob=/opt/reagent/reagent-*
Unit=reagent-upgrade.service

[Install]
#WantedBy=sysinit.target
